# zxtouchd/layout/DEBIAN/postinst
#!/bin/sh
launchctl load /Library/LaunchDaemons/com.zjx.zxtouchd.plist
