# zxtouchd/layout/DEBIAN/postinst
#!/bin/sh
launchctl load /var/jb/Library/LaunchDaemons/com.zjx.zxtouchd.plist
